# dunst/.config/systemd/user/dunst.service
[Unit]
Description=Dunst notification daemon

[Service]
Type=simple
ExecStart=/usr/bin/dunst

[Install]
WantedBy=default.target
